# rootfs/etc/rc.local
#!/bin/sh
################################################################################
#	Copyright (C) 2019		Alejandro Colomar Andrés	       #
#	SPDX-License-Identifier:	ISC				       #
################################################################################
# Deprecated in systemd systems!
# Please use systemd version
################################################################################

/etc/rc.local.d/random-mac
